# ss1/system/etc/profile.d/logout.sh
on_logout() {
    [ -f "$HISTFILE" ] && rm "$HISTFILE"
    HISTFILE=
}

trap 'on_logout' EXIT
